# bin/profile.d/mount_env.sh
#!/usr/bin/env bash
sudo ln -s /home/douglasmg7/code/webserver/bin/profile.d/webserver.sh /etc/profile.d/webserver.sh
echo System must be restarted, to load environment variables
